# serverstats.cron
#
# cron job for the serverstats package
#
* * * * *	stats	/usr/share/serverstats/update.php
